# output/lib/systemd/system/ethercat.service
#
# EtherCAT Master Kernel Modules
#

[Unit]
Description=EtherCAT Master Kernel Modules

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/wzz/ethercat/ethercat-1.5.2/output/sbin/ethercatctl start
ExecStop=/home/wzz/ethercat/ethercat-1.5.2/output/sbin/ethercatctl stop

[Install]
WantedBy=multi-user.target
